# solvia-metrics.service
[Unit]
Description=Solvia Blockchain Metrics Exporter
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/solvia
ExecStart=/usr/bin/node /root/solvia/src/monitoring/metrics-exporter.js
Environment=NODE_ENV=production
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
